# arcade/scripts/quartermaster-audit.timer
[Unit]
Description=Run Quartermaster audit weekly

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
