# systemd/halo-caddy.service
[Unit]
Description=halo-ai Caddy Reverse Proxy
After=network.target

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
ExecStart=/usr/local/bin/caddy run --config /srv/ai/configs/Caddyfile
ExecReload=/usr/local/bin/caddy reload --config /srv/ai/configs/Caddyfile
Restart=on-failure
RestartSec=5
AmbientCapabilities=CAP_NET_BIND_SERVICE
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai %h/.local/share/caddy

[Install]
WantedBy=multi-user.target
