# systemd/halo-n8n.service
[Unit]
Description=n8n Workflow Automation
After=network.target

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
WorkingDirectory=/srv/ai/n8n
Environment=N8N_PORT=5678
Environment=N8N_LISTEN_ADDRESS=127.0.0.1
Environment=N8N_USER_FOLDER=/srv/ai/n8n/data
ExecStart=/usr/local/bin/node /srv/ai/n8n/packages/cli/bin/n8n start
Restart=on-failure
RestartSec=5
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai

[Install]
WantedBy=multi-user.target
