# systemd/halo-llama-server.service
[Unit]
Description=llama.cpp Server (HIP/ROCm)
After=network.target

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
EnvironmentFile=/srv/ai/configs/rocm.env
Environment=ROCBLAS_USE_HIPBLASLT=1
ExecStart=/srv/ai/llama-cpp/build-vulkan/bin/llama-server     --host 127.0.0.1 --port 8081     --no-mmap -ngl 99 -fa on     --model /srv/ai/models/qwen3-30b-a3b-q4_k_m.gguf
Restart=on-failure
RestartSec=5
SupplementaryGroups=video render
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai

[Install]
WantedBy=multi-user.target
