# systemd/halo-vane.service
[Unit]
Description=Vane (Perplexica) Deep Research
After=network.target halo-searxng.service halo-lemonade.service

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
WorkingDirectory=/srv/ai/vane
Environment=PORT=3001
Environment=HOSTNAME=127.0.0.1
Environment=SEARXNG_API_URL=http://localhost:8888
ExecStart=/usr/local/bin/node /srv/ai/vane/.next/standalone/server.js
Restart=on-failure
RestartSec=5
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai

[Install]
WantedBy=multi-user.target
