# systemd/halo-security-audit-weekly.timer
[Unit]
Description=halo-ai Weekly Security Audit Archive Timer

[Timer]
OnCalendar=Mon *-*-* 06:30:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
