# systemd/halo-whisper-server.service
[Unit]
Description=whisper.cpp Server
After=network.target

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
EnvironmentFile=/srv/ai/configs/rocm.env
ExecStart=/srv/ai/whisper-cpp/build/bin/whisper-server     --host 127.0.0.1 --port 8082     --model /srv/ai/models/whisper-large-v3-turbo.bin
Restart=on-failure
RestartSec=5
SupplementaryGroups=video render
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai

[Install]
WantedBy=multi-user.target
