# systemd/halo-backup.service
[Unit]
Description=halo-ai Nightly Backup

[Service]
Type=oneshot
User=<YOUR_USER>
ExecStart=/srv/ai/scripts/halo-backup.sh /srv/ai/backups
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai
