# systemd/halo-open-webui.service
[Unit]
Description=Open WebUI
After=network.target halo-lemonade.service

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
WorkingDirectory=/srv/ai/open-webui
Environment=DATA_DIR=/srv/ai/open-webui/data
Environment=OLLAMA_BASE_URL=http://localhost:8080
ExecStart=/srv/ai/open-webui/.venv/bin/open-webui serve --host 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=5
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai

[Install]
WantedBy=multi-user.target
