# systemd/halo-backup.timer
[Unit]
Description=Daily Halo AI backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
